# crontab
*/5 * * * * /usr/local/bin/python /app/call_api.py
